# roles/Gunicorn/templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
SocketUser={{ os_user }}
SocketGroup=www-data
SocketMode=0660

[Install]
WantedBy=sockets.target
